# wine/autostart/wine.desktop
[Desktop Entry]
Name=Wine
Type=Application
Icon=wine
Exec=wineboot

